# profile.d/macvim.sh
if [[ -d /Applications/MacVim.app/Contents/bin ]]; then
    pathmunge /Applications/MacVim.app/Contents/bin
fi
