# profile.d/oracle.sh
# oracle settings

if [[ -e /opt/oracle-instantclient/current ]]; then
   pathmunge /opt/oracle-instantclient/current after
fi
