# profile.d/oracle.sh
# oracle settings

if [[ -e /opt/oracle-instantclient/current ]]; then
   PATH="$PATH:/opt/oracle-instantclient/current"
fi
