# profile.d/orbstack.sh
if [[ -e ~/.orbstack/shell/init.bash ]]; then
    source ~/.orbstack/shell/init.bash
fi
