# profile.d/python.sh
#!/bin/bash

export VIRTUAL_ENV_DISABLE_PROMPT=y
[[ -f ~/.python/current/bin/activate ]] && source ~/.python/current/bin/activate

