# profile.d/python.sh
export VIRTUAL_ENV_DISABLE_PROMPT=y
[[ -f ~/.local/bin/activate ]] && source ~/.local/bin/activate
