# profile.d/saml.sh
#!/usr/bin/env bash

downloadInCommonMetadata() {
    wget -O incommon-metadata.xml https://mdq.incommon.org/entities
}
